# deploy/miraeboy.service
[Unit]
Description=miraeboy - Conan v2 Package Server
Documentation=https://github.com/corbie79/miraeboy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=miraeboy
Group=miraeboy

WorkingDirectory=/etc/miraeboy
ExecStart=/usr/local/bin/miraeboy
ExecReload=/bin/kill -HUP $MAINPID

Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s

# 보안 강화
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/miraeboy
PrivateTmp=true
PrivateDevices=true

# 로그는 journald로 수집
StandardOutput=journal
StandardError=journal
SyslogIdentifier=miraeboy

[Install]
WantedBy=multi-user.target
